# lesson-07/myservice/myservice.service
[Unit]
Description=myservice - my log scanner

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/myservice.conf
ExecStart=/usr/bin/myservice.sh